# autologout/src/deb/control/postinst
#!/bin/sh
chmod 755 /etc/init.d/autologout
update-rc.d -f autologout defaults