# autologout/src/deb/control/postinst
#!/bin/sh
update-rc.d autologout defaults