# autologout/src/deb/control/postinst
#!/bin/sh
update-rc.d -f autologout defaults